# com.henrybravo.mlx-router.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.henrybravo.mlx-router</string>
    <key>ProgramArguments</key>
    <array>
        <string>__INSTALL_DIR__/venv/bin/python</string>
        <string>__INSTALL_DIR__/main.py</string>
        <string>--config</string>
        <string>__INSTALL_DIR__/config.json</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__INSTALL_DIR__</string>
    <key>StandardOutPath</key>
    <string>__INSTALL_DIR__/logs/mlx_router.log</string>
    <key>StandardErrorPath</key>
    <string>__INSTALL_DIR__/logs/mlx_router.error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HF_HUB_CACHE</key>
        <string>__HOME__/models</string>
    </dict>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>